# service/mlb-stats-api-data.service.d/variables.conf
[Service]
Environment=MONGO_DB_NAME=mlb-stats
Environment=MONGO_IP=localhost
Environment=MONGO_PASSWORD=