# service/mlb-stats-api-data.service
[Unit]
Description=Nightly job to collect MLB stats data
After=network.target
[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/mlb-stats-api
ExecStart=/home/ubuntu/mlb-stats-api/venv/bin/python data_capture.py
Restart=always
[Install]
WantedBy=multi-user.target